# private_dot_config/systemd/user/agent-deck-conductor-bridge.service.d/override.conf
# Drop-in override for agent-deck's GENERATED unit
# (~/.config/systemd/user/agent-deck-conductor-bridge.service, written by
# `conductor setup`). Lives in a .d/ dir so it survives agent-deck regenerating
# the main .service file. systemd applies drop-ins last, so this PATH wins.
#
# Why (2026-06-09): Nix paths must come FIRST so `tmux` resolves to the Nix 3.6a
# binary that matches the running agent-deck tmux server. With system tmux 3.4
# first, the bridge mis-detected the live conductor as "not running", spawned a
# phantom session and dropped the Telegram relay. /opt/homebrew dropped (macOS).
#
# EnvironmentFile (additive in a drop-in) sources the Telegram bot token
# (TELEGRAM_BOT_TOKEN) from a 0600 .env, resolved by bridge.py's _resolve_secret
# via config.toml token = "$TELEGRAM_BOT_TOKEN". Leading "-" makes it optional.
# The .env stays untracked/secret (org policy: never commit .env).
[Service]
Environment=PATH=/home/paper/.nix-profile/bin:/nix/var/nix/profiles/default/bin:/home/paper/.local/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin
EnvironmentFile=-/home/paper/.sbx/secrets.env
